# bitcoin.service
[Unit]
Description=Bitcoin daemon
Documentation=man:bitcoind(1)
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/bitcoind -conf=${CONFIG_FILE} -datadir=${DATA_DIR} -pid=${PID_FILE} -debuglogfile=${DEBUG_FILE} $OPTIONS
EnvironmentFile=/etc/sysconfig/bitcoin
User=bitcoin
Group=bitcoin

Restart=on-failure
TimeoutStopSec=120
TimeoutStartSec=60
StartLimitInterval=240
StartLimitBurst=5

PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
